# deploy/systemd/user-service.service
[Unit]Description=User Service - Discord OAuth2 AuthenticationAfter=network.target postgresql.serviceWants=postgresql.service[Service]Type=simpleUser=deployGroup=deployWorkingDirectory=/opt/user-serviceExecStart=/opt/user-service/user-serviceRestart=alwaysRestartSec=5StandardOutput=journalStandardError=journalSyslogIdentifier=user-service# Load environment variables from .env fileEnvironmentFile=/opt/user-service/.env# Security settingsNoNewPrivileges=truePrivateTmp=trueProtectSystem=strictProtectHome=trueReadWritePaths=/opt/user-serviceCapabilityBoundingSet=CAP_NET_BIND_SERVICE# Resource limitsLimitNOFILE=65536LimitNPROC=4096[Install]WantedBy=multi-user.target